# launchd/com.zoraasi.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.zoraasi.daemon</string>
  <key>ProgramArguments</key><array>
    <string>/opt/homebrew/anaconda3/bin/python3.13</string>
    <string>/Users/christophermichaelbaird/Documents/New OpenCode Project/zoraos/scripts/zora_daemon.py</string>
    <string>--goal</string><string>Quiet overnight reflection and journaling. Keep entries warm, restful, and concise for Christopher.</string>
    <string>--provider</string><string>ollama</string>
    <string>--model</string><string>glm4:9b</string>
    <string>--continuous</string><string>--interval-minutes</string><string>480</string>
    <string>--daily-tasks</string><string>3</string>
    <string>--daily-tool-calls</string><string>12</string>
    <string>--daily-tokens</string><string>32000</string>
    <string>--task-iterations</string><string>10</string>
    <string>--task-tokens</string><string>16000</string>
    <string>--task-wall-seconds</string><string>300</string>
    <string>--no-ws</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/christophermichaelbaird/Documents/New OpenCode Project/zoraos</string>
  <key>KeepAlive</key><true/><key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/christophermichaelbaird/Documents/New OpenCode Project/zoraos/data/daemon/launchd.out.log</string>
  <key>StandardErrorPath</key><string>/Users/christophermichaelbaird/Documents/New OpenCode Project/zoraos/data/daemon/launchd.err.log</string>
</dict></plist>
