# launchd/com.zoraasi.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.zoraasi.daemon</string>
  <key>ProgramArguments</key><array>
    <string>/opt/homebrew/anaconda3/bin/python3.13</string>
    <string>/Users/christophermichaelbaird/Documents/New OpenCode Project/zoraos/scripts/zora_daemon.py</string>
    <string>--goal</string><string>Produce a concise local operational brief from the facts supplied in this task. State one verified system condition, one unresolved risk, and one safe next action. Do not use tools, browse, write files, or claim capabilities beyond the supplied evidence.</string>
    <string>--provider</string><string>openrouter</string>
    <string>--model</string><string>google/gemma-4-26b-a4b-it:free</string>
    <string>--continuous</string><string>--interval-minutes</string><string>30</string>
    <string>--no-memory</string><string>--daily-tasks</string><string>6</string>
    <string>--daily-tool-calls</string><string>0</string><string>--daily-tokens</string><string>12000</string>
    <string>--task-tool-calls</string><string>0</string><string>--task-iterations</string><string>3</string>
    <string>--task-tokens</string><string>1200</string><string>--task-wall-seconds</string><string>180</string>
    <string>--ws-port</string><string>8766</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/christophermichaelbaird/Documents/New OpenCode Project/zoraos</string>
  <key>KeepAlive</key><true/><key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/christophermichaelbaird/Documents/New OpenCode Project/zoraos/data/daemon/launchd.out.log</string>
  <key>StandardErrorPath</key><string>/Users/christophermichaelbaird/Documents/New OpenCode Project/zoraos/data/daemon/launchd.err.log</string>
</dict></plist>
